# systemd/dhconsole.service
[Unit]
Description=DigiHub Web Console (ttyd + login)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/dhconsoled __DIGIHUB_USER__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
