# systemd/dhjs8call.service
[Unit]
Description=DigiHub JS8Call (VNC + noVNC remote desktop)
After=network.target sound.target

[Service]
Type=simple
User=__DIGIHUB_USER__
SupplementaryGroups=dialout audio
ExecStart=/usr/local/bin/dhjs8calld
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
